# stow/aconfmgr/.config/aconfmgr/files/etc/systemd/system/gdm.service.d/override.conf
[Service]
ExecStartPre=/usr/sbin/gdm-monitors-update
